# content/opt/mysql.conf
[program:mysql]
command=/bin/bash -c "/opt/supervisor/mysql_supervisor"
redirect_stderr=true
autorestart=true
