# content/opt/mysql.conf
[program:mysql]
command=/bin/bash -c "/opt/supervisor/mysql_supervisor"
redirect_stderr=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
