# rundeck.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/bin/bash -c "/opt/supervisor/mysql_supervisor"
redirect_stderr=true

[program:rundeck]
command=/opt/supervisor/rundeck
redirect_stderr=true
